# conf/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/log/supervisor/superfisord.pid
[program:Apache]
command=httpd -DFOREGROUND
process_name=httpd
